# etc/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:nginx]
command=/usr/sbin/nginx

[program:rally-ci]
command=/usr/local/bin/rally-ci /home/rally/config.yaml
